# da/workflows/server/crontab/crontab.add
# start monitor daemon
* * * * * . ~/.llview_server_rc ; perl "$LLVIEW_HOME/da/workflows/server/crontab/serverAll.pl"

