# da/workflows/remote/crontab/crontab.add
* * * * * . ~/.llview_remote_rc ; perl "$LLVIEW_HOME/da/workflows/remote/crontab/remoteAll.pl"
